# stowaway/root/etc/s6-overlay/s6-rc.d/init-wireguard-module/run
#!/usr/bin/with-contenv bash

mkdir -p /config/{templates,coredns}

echo "Uname info: $(uname -a)"
