# cargo/root/etc/s6-overlay/s6-rc.d/svc-wireguard/finish
#!/usr/bin/with-contenv bash

exec /bin/bash /usr/bin/wg-quick down /config/wg0.conf
