# cargo/root/etc/s6-overlay/s6-rc.d/svc-wireguard/finish
#!/usr/bin/with-contenv bash

wg-quick down /config/wg0.conf
