# cargo/root/etc/s6-overlay/s6-rc.d/svc-wireguard/run
#!/usr/bin/with-contenv bash

exec /bin/bash /usr/bin/wg-quick up /config/wg0.conf
