# stowaway/root/etc/s6-overlay/s6-rc.d/svc-tcp2udp/run
#!/usr/bin/with-contenv bash
 
tcp2udp --tcp-listen 0.0.0.0:51821 --udp-forward 127.0.0.1:51820