# stowaway/root/etc/services.d/wireguard/run
#!/usr/bin/with-contenv bash

_term() {
  echo "Caught SIGTERM signal!"
  wg-quick down wg0
}

trap _term SIGTERM

wg-quick up /etc/wireguard/wg0.conf

sleep infinity &

wait
