# stowaway/old/services.d/rsync/run
#!/usr/bin/with-contenv bash


if [ ! -f /var/run/rsyncd.pid ]; then
  exec rsync --daemon
else
  sleep infinity
fi