# carrier2/etc/inittab
::once:carrier2 -c /tmp/config.yaml -d
::respawn:tail -F /tmp/carrier.log