# stowaway/old/services.d/nginx/run
#!/usr/bin/with-contenv bash


if [ ! -f /tmp/nginx.pid ]; then
  exec nginx -g "daemon off;"
else
  sleep infinity
fi