# stowaway/root/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/usr/bin/with-contenv bash

s6-echo "Starting Nginx service"
nginx -g "daemon off;"
