# cargo/root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

if [ ! -f /var/run/crond.pid  ]; then
  exec cron
else
  sleep infinity
fi